# elements/cc-common/static/etc/systemd/system/cache-vendordata.timer
[Unit]
Description=Run vendordata fetch every hour

[Timer]
OnBootSec=30s
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
